# OpenWRT/overlay/upper/etc/hotplug.d/iface/60-qos
#!/bin/sh
logger -t hotplug "Start 60-qos for $ACTION on $INTERFACE"
[ "$ACTION" = ifup ] && /etc/init.d/qos enabled && /usr/lib/qos/generate.sh interface "$INTERFACE" | sh
