# AsusWRT-Merlin/jffs/opt/lib/opkg/info/ntpd.list
/opt/etc/ntp.conf
/opt/sbin/ntpd
/opt/etc/init.d/S77ntpd
/opt/etc/hotplug.d/iface/20-ntpd
