# jffs/opt/lib/opkg/info/pixelserv-tls.list
/opt/bin/pixelserv-tls
/opt/etc/init.d/S80pixelserv-tls
