# OpenWRT/overlay/upper/usr/lib/opkg/info/darkstat.list
/etc/init.d/darkstat
/etc/config/darkstat
/usr/sbin/darkstat
