# OpenWRT/overlay/upper/usr/lib/opkg/info/hostapd-common.list
/lib/netifd/hostapd.sh
/etc/init.d/wpad
/etc/rc.button/wps
