# OpenWRT/overlay/upper/usr/lib/opkg/info/ksmbd-server.list
/etc/init.d/ksmbd
/etc/ksmbd/smb.conf.template
/etc/ksmbd/configuration.txt
/usr/sbin/ksmbd.mountd
/etc/config/ksmbd
/etc/ksmbd/ksmbd.config.example
/lib/upgrade/keep.d/ksmbd-server
