# OpenWRT/overlay/upper/usr/lib/opkg/info/rssileds.list
/usr/libexec/led-trigger/rssi
/etc/init.d/rssileds
/usr/sbin/rssileds
